# scripts/hydro.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
HOME=/

0 1 * * * root /var/www/FAHS/scripts/hydro_killandstart.sh >/dev/null 2>&1
0 6,12,18 * * * root /var/www/FAHS/scripts/take_picture.sh >/dev/null 2>&1
